# containers/systemd/docling/docling.container
# Generated by Chezmoi - DO NOT EDIT MANUALLY

[Unit]
Description=Docling - Document processing service
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/ds4sd/docling:latest
AutoUpdate=registry
ContainerName=docling

# Network
Network=llm.network

# Published port for Caddy to access (optional, can be internal only)
PublishPort={{ .published_ports.docling }}:{{ .ports.docling }}

# Health check
HealthCmd="curl -f http://localhost:{{ .ports.docling }}/health || exit 1"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
