# containers/systemd/openwebui-postgres/openwebui-postgres.container
# Generated by Chezmoi - DO NOT EDIT MANUALLY

[Unit]
Description=PostgreSQL database for Open WebUI
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/pgvector/pgvector:pg16
AutoUpdate=registry
ContainerName=openwebui-postgres

# Network - NO published port, only internal access
Network=llm.network

# Environment
Environment=POSTGRES_DB=openwebui
Environment=POSTGRES_USER={{ .database.openwebui_postgres_user }}
Environment=POSTGRES_PASSWORD={{ .database.openwebui_postgres_password }}

# Volume
Volume=openwebui-postgres.volume:/var/lib/postgresql/data:Z

# Health check
HealthCmd="pg_isready -U {{ .database.openwebui_postgres_user }}"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
