# containers/systemd/searxng/searxng.container
# Generated by Chezmoi - DO NOT EDIT MANUALLY

[Unit]
Description=SearXNG - Privacy-respecting metasearch engine
After=network-online.target llm.network.service
After=searxng-redis.service
Requires=llm.network.service
Wants=searxng-redis.service

[Container]
Image=docker.io/searxng/searxng:latest
AutoUpdate=registry
ContainerName=searxng

# Network
Network=llm.network

# Published port for Caddy to access
PublishPort={{ .published_ports.searxng }}:{{ .ports.searxng }}

# Volume for settings
Volume=searxng.volume:/etc/searxng:Z

# Environment
Environment=INSTANCE_NAME=searxng
Environment=BASE_URL=http://0.0.0.0:{{ .ports.searxng }}
Environment=REDIS_URL=redis://searxng-redis:{{ .ports.redis }}/0

# Health check
HealthCmd="wget --no-verbose --tries=1 --spider http://localhost:{{ .ports.searxng }} || exit 1"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
