# containers/systemd/sst/sst.container
[Unit]
Description=Faster Whisper STT (Speech-to-Text) Service
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/fedirz/faster-whisper-server:latest-cpu
AutoUpdate=registry
ContainerName=sst

# Network
Network=llm

# Published port for Caddy access
PublishPort=127.0.0.1:8000:8000

# Volume for model cache
Volume=sst.volume:/root/.cache:Z

# Environment
Environment=WHISPER_MODEL=base
Environment=WHISPER_LANGUAGE=en

# Health check
HealthCmd="wget --no-verbose --tries=1 --spider http://localhost:8000/health || exit 1"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
