# containers/systemd/litellm/litellm.container
# Generated by Chezmoi - DO NOT EDIT MANUALLY

[Unit]
Description=LiteLLM Proxy Service
After=network-online.target llm.network.service
After=litellm-postgres.service litellm-redis.service
Requires=llm.network.service
Wants=litellm-postgres.service litellm-redis.service

[Container]
Image=ghcr.io/berriai/litellm:main-stable
AutoUpdate=registry
ContainerName=litellm

# Network
Network=llm.network

# Published port for Caddy to access
PublishPort={{ .published_ports.litellm }}:{{ .ports.litellm }}

# Mount config file
Volume=%h/.config/containers/systemd/litellm/litellm.yaml:/app/config.yaml:ro,Z

# Environment variables
Environment=LITELLM_MASTER_KEY={{ .api_keys.litellm_master }}
Environment=DATABASE_URL=postgresql://{{ .database.litellm_postgres_user }}:{{ .database.litellm_postgres_password }}@litellm-postgres:{{ .ports.postgres }}/litellm
Environment=REDIS_HOST=litellm-redis
Environment=REDIS_PORT={{ .ports.redis }}
Environment=STORE_MODEL_IN_DB=true
Environment=USE_AIOHTTP_TRANSPORT=true

# API Keys
Environment=OPENAI_API_KEY={{ .api_keys.openai }}
Environment=ANTHROPIC_API_KEY={{ .api_keys.anthropic }}
Environment=GEMINI_API_KEY={{ .api_keys.gemini }}

# Execution
Exec=--config /app/config.yaml --port {{ .ports.litellm }} --num_workers 2

# Health check
HealthCmd="curl -f http://localhost:{{ .ports.litellm }}/health/readiness || exit 1"
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=30s

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
