# containers/systemd/litellm/litellm-postgres.container
# Generated by Chezmoi - DO NOT EDIT MANUALLY

[Unit]
Description=PostgreSQL database for LiteLLM
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/postgres:16
AutoUpdate=registry
ContainerName=litellm-postgres

# Network - NO published port, only internal access
Network=llm.network

# Environment
Environment=POSTGRES_DB=litellm
Environment=POSTGRES_USER={{ .database.litellm_postgres_user }}
Environment=POSTGRES_PASSWORD={{ .database.litellm_postgres_password }}

# Volume
Volume=litellm-postgres.volume:/var/lib/postgresql/data:Z

# Health check
HealthCmd="pg_isready -U {{ .database.litellm_postgres_user }}"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=5

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
