# containers/systemd/litellm/litellm-redis.container
# Generated by Chezmoi - DO NOT EDIT MANUALLY

[Unit]
Description=Redis cache for LiteLLM
After=network-online.target llm.network.service
Requires=llm.network.service

[Container]
Image=docker.io/redis:latest
AutoUpdate=registry
ContainerName=litellm-redis

# Network - NO published port, only internal access
Network=llm.network

# Volume
Volume=litellm-redis.volume:/data:Z

# Health check
HealthCmd="redis-cli ping || exit 1"
HealthInterval=30s
HealthTimeout=5s
HealthRetries=3
HealthStartPeriod=10s

[Service]
Slice=llm.slice
TimeoutStartSec=900
Restart=on-failure
RestartSec=10

[Install]
WantedBy=scroll-session.target
